# extra/supervisor.example.conf
[program:cares]
command=/path/to/site/cares --database 'dbname=cares user=cares' --port 8080
directory=/path/to/site
user=nobody
autostart=true
autorestart=true
redirect_stderr=True
